# Ubuntu/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn run:app -b localhost:6000
directory=/deploy/app
autorestart  = true
autostart    = true
startsecs    = 3
startretires = 3